# recipes-core/companion/files/mediamtx.service
[Unit]
Description=MediaMTX RTSP/WHEP streaming server
After=network.target

[Service]
Type=simple
User=root
# mediamtx's own default config search (relative to CWD: ./mediamtx.yml,
# ./rtsp-simple-server.yml) never looks at /etc — confirmed live on-device
# via "configuration file not found ... using an empty configuration",
# which silently discarded every custom setting in mediamtx.yml (WHEP
# port, RTSP port, etc.) and fell back to defaults, one of which conflicted
# on :8000 and crash-looped the service. Pass the real path explicitly.
ExecStart=/usr/bin/mediamtx /etc/mediamtx.yml
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
